# pgl/pglcmd/init/pgl.service.in
[Unit]
Description=PeerGuardian Linux - an IP Blocker
Documentation=man:pgld(1) file://@data_root_dir@/doc/pgl/README.blocklists
After=network.target
ConditionPathExists=|@CONF_DIR@/pgl/blocklists.list
ConditionDirectoryNotEmpty=|@LIB_DIR@/pgl

[Service]
BusName=org.netfilter.pgl
ExecStart=@BIN_DIR@/pglcmd start
ExecStop=@BIN_DIR@/pglcmd stop
PIDFile=@PID_DIR@/pgld.pid

[Install]
WantedBy=multi-user.target
