# pgl/pglcmd/init/pgl.service.in
[Unit]
Description=PeerGuardian Linux - an IP Blocker
Documentation=man:pgld(1)
Documentation=file://@data_root_dir@/doc/pgl/README.blocklists
After=network.target syslog.target
After=firehol.service firestarter.service firewalld.service ufw.service

[Service]
Type=forking
ExecStart=@BIN_DIR@/pglcmd start
ExecStop=@BIN_DIR@/pglcmd stop
PIDFile=@PID_DIR@/pgld.pid

[Install]
WantedBy=multi-user.target
