# pgl/pglcmd/debian/preinst
#! /bin/sh
# preinst script for pglcmd
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#		* <new-preinst> `install'
#		* <new-preinst> `install' <old-version>
#		* <new-preinst> `upgrade' <old-version>
#		* <old-preinst> `abort-upgrade' <new-version>
#
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
	install|upgrade)
		# Remove other init files to prevent conflict with pglcmd by several identical services.
		for i in blockcontrol moblock moblock-nfq moblock-ipq peerguardnf peerguardian ; do
			if [ -x "/etc/init.d/$i" ] ; then
				rm /etc/init.d/$i || exit $?
				update-rc.d $i remove >/dev/null || exit $?
			fi
		done
		# Remove old system files.
        [ ! -f /etc/logrotate.d/moblock-control ] || rm /etc/logrotate.d/moblock-control
        [ ! -f /etc/network/if-up.d/moblock ] || rm /etc/network/if-up.d/moblock
        [ ! -f /etc/cron.daily/moblock ] || rm /etc/cron.daily/moblock
        [ ! -f /etc/logrotate.d/blockcontrol ] || rm /etc/logrotate.d/blockcontrol
        [ ! -f /etc/network/if-up.d/blockcontrol ] || rm /etc/network/if-up.d/blockcontrol
        [ ! -f /etc/cron.daily/blockcontrol ] || rm /etc/cron.daily/blockcontrol
		;;
	abort-upgrade)
		;;
	*)
		echo "preinst called with unknown argument \`$1'" >&2
		exit 1
		;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
