# pgl/pgld/debian/cron.weekly
#!/bin/sh

ENABLED=0
BLOCKLISTPRO_NAMES=

test -f /etc/default/pgl && . /etc/default/pgl

if [ "$ENABLED" = "0" -o -z "$BLOCKLISTPRO_NAMES" ]; then
    exit 0
fi

/usr/lib/pgl/dl-blocklistpro.pl $BLOCKLISTPRO_NAMES
