# pgl/pgld/debian/cron.daily
#!/bin/sh

ENABLED=0
BLOCKLIST_URL=

test -f /etc/default/pgl && . /etc/default/pgl

if [ "$ENABLED" = "0" -o -z "$BLOCKLIST_URL" ]; then
    exit 0
fi

wget -q -P/var/lib/pgl/ -N $BLOCKLIST_URL
